# systemctl-win/tests/service_units/i.service
[Unit]
Description=Test-Service-I
Wants=a.service f.service

[Service]
ExecStart=/bin/true
