# systemctl-win/tests/service_units/test-d.service
[Unit]
Description=Test D:Cyclic
After=test-b.service
Before=test-a.service
Requires=test-a.service

[Service]
ExecStart=/bin/true
