# systemctl-win/tests/service_units/test-a.service
[Unit]
Description=Test A
Requires=test-b.service
Before=test-b.service

[Service]
ExecStart=/bin/true
