# systemctl-win/tests/service_units/test-h.service
[Unit]
Description=Test-Service-H
Wants=test-g.service

[Service]
ExecStart=/bin/true
