# systemctl-win/tests/service_units/test-i.service
[Unit]
Description=Test-Service-I
Wants=test-a.service test-f.service

[Service]
ExecStart=/bin/true
