# systemctl-win/tests/service_units/test-c.service
[Unit]
Description=Test C
Requires=test-a.service

[Service]
ExecStart=/bin/true
