# systemctl-win/tests/service_units/test-g.service
[Unit]
Description=G
Conflicts=test-e.service

[Service]
ExecStart=/bin/true
