# systemctl-win/tests/service_units/test-b.service
[Unit]
Description=Test B
Wants=test-f.service

[Service]
ExecStart=/bin/true
