# systemctl-win/tests/service_units/test-e.service
[Unit]
Description=Test E:Cyclic
After=test-b.service
Before=test-a.service
Wants=test-a.service

[Service]
ExecStart=/bin/true
